# target/linux/vf610/base-files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
ttyLP0::askfirst:/bin/ash --login
